# override/etc/services.d/00-crond/run
#!/usr/bin/with-contenv sh

/etc/periodic/15min/watcher-crond
exec /usr/sbin/crond -L /mnt/config/log/crond.log -l 6 -f;
