# override/etc/s6-overlay/s6-rc.d/logrotate_setup/up
#!/command/execlineb
# SPDX-FileCopyrightText: 2022 - 2024 Daniel Wolf <nephatrine@gmail.com>
#
# SPDX-License-Identifier: ISC

foreground {
  if -n { /usr/bin/test -d /mnt/config/etc/logrotate.d }
  /command/s6-setuidgid guardian /bin/mkdir -p /mnt/config/etc/logrotate.d
}
foreground {
  if -n { /usr/bin/test -f /mnt/config/etc/logrotate.conf }
  /command/s6-setuidgid guardian /bin/cp /etc/logrotate.conf /mnt/config/etc/
}
foreground { /command/s6-setuidgid guardian /bin/cp -Rn /etc/logrotate.d /mnt/config/etc/ }
